# liquid-stats.timer
[Unit]
Description=Run stats every 2 minutes, 1 minute after boot

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
